# setup/com.anki.stats-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anki.stats-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/lori/Workspace/VSCodeProjects/anki-stats-sync/src/sync.py</string>
        <string>--push</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run daily at 23:00 -->
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/anki-stats-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/anki-stats-sync.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/lori/Workspace/VSCodeProjects/anki-stats-sync</string>
</dict>
</plist>
